# setup/squid.conf
[program:squid]
user = nobody
command = /usr/bin/squid -N
umask = 002
